# system_files/usr/lib/systemd/user/gamescope-session-plus@ogui-steam.service.d/10-apps.conf
# Drop-in override for gamescope-session-plus to launch headless apps.
# This does not modify the upstream session; it only adds a dependency.

[Unit]
# Pull in the app target without making app failures break the Gamescope session.
Wants=gamescope-apps.target
