# system_files/usr/lib/systemd/user/gamescope-session-plus@.service.d/10-apps.conf
# Drop-in override for gamescope-session-plus to launch headless apps
# This does NOT modify the upstream session - it only adds a dependency

[Unit]
# Pull in headless apps when gamescope session starts
# Uses Wants= not Requires= so apps failure doesn't break the session
Wants=gamescopeApps.service
